# tags/release-0.2/rc/halvm
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Foundation
# Distributed under the terms of the GNU General Public License, v2 or later
# $Header: $

depend() {
	need hald
	use logger
}

start() {
	ebegin "Starting HalVM"
	start-stop-daemon --start --pidfile /var/run/halvm.pid --exec /usr/sbin/halvm > /dev/null
	eend $?
}

stop() {
	ebegin "Stopping HalVM"
	start-stop-daemon --stop --pidfile /var/run/halvm.pid
	eend $?
}
